# backend/therapydash-backend.service
# /etc/systemd/system/therapydash-backend.service
[Unit]
Description=TherapyDash Express Backend
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/therapydash-demo/backend
ExecStart=/usr/bin/node app.js
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=therapydash-backend
Environment=NODE_ENV=production
Environment=PORT=3000

[Install]
WantedBy=multi-user.target